# dot_config/systemd/user/vault-sync.service
[Unit]
Description=Sync Obsidian vault to git remote

[Service]
Type=oneshot
ExecStart=%h/.local/bin/vault-sync
